# ots/DEBIAN/postinst
#!/bin/bash

systemctl enable uwsgi

exit 0
