# ots/DEBIAN/prerm
#!/bin/bash

systemctl stop uwsgi
systemctl disable uwsgi

exit 0
